# supervisor_smb.conf
[program:smbd]
command=/usr/sbin/smbd --foreground
autorestart=true
